# phobos-panel.service
[Unit]
Description=Phobos VPN Web Panel
After=network.target wg-quick@wg0.service
Wants=wg-quick@wg0.service

[Service]
Type=simple
WorkingDirectory=/opt/phobos-panel
ExecStart=/usr/bin/gunicorn -w 1 -b 0.0.0.0:8443 app:app
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
